# shells/systemd/surreal-8020-rocksdb.service
[Unit]
Description=SurrealDB plant (8020, rocksdb, same as local ams-8020.db)
Documentation=https://surrealdb.com/docs
After=network-online.target
Wants=network-online.target
# 与 surrealkv 的 surrealdb.service 互斥，避免争用 8020
Conflicts=surrealdb.service

[Service]
Type=simple
WorkingDirectory=/root/surreal_data
Environment=REMOTE_SURREAL_DATA_DIR=/root/surreal_data
Environment=SURREAL_DB_DIR_NAME=ams-8020.db
Environment=SURREAL_USER=root
Environment=SURREAL_PASS=root
Environment=SURREAL_BIND=0.0.0.0:8020
# Ubuntu 常见安装路径为 /usr/local/bin/surreal；也可用 EnvironmentFile 覆盖
Environment=SURREAL_BIN=/usr/local/bin/surreal
ExecStart=/root/shells/run_surreal_8020_ubuntu.sh
Restart=on-failure
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
